# test-xooj/init_data/check_docker/checker.conf
[program:checker_server]
command=python server.py
directory=/home/rpc
autostart=true
autorestart=true
user=root